# lesson6/mqtt-monitor.service
[Unit]
Description=MQTT Sensor Monitor Web Application
After=network.target mosquitto.service
Wants=mosquitto.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/Documents/GitHub/20205_10_26_chihlee_pi_pico_Steve/lesson6
ExecStart=/home/pi/.local/bin/uv run python app_flask.py
Environment=PATH=/home/pi/.local/bin:/usr/local/bin:/usr/bin:/bin
Restart=always
RestartSec=5
StandardOutput=append:/home/pi/Documents/GitHub/20205_10_26_chihlee_pi_pico_Steve/lesson6/app.log
StandardError=append:/home/pi/Documents/GitHub/20205_10_26_chihlee_pi_pico_Steve/lesson6/app.log

[Install]
WantedBy=multi-user.target

